# recipes-graphics/linux-input/linux-input/12-touch-inputattach.rules
ACTION=="add", KERNEL=="ttymxc1", RUN+="/usr/sbin/inputattach --daemon --embtom /dev/%k"
